# debian/budgie-sddm-theme.sddm-backgrounds-cache.path
[Unit]
Description=Watch AccountsService user files for changes
Requires=sddm-backgrounds-cache.service

[Path]
PathModified=/var/lib/AccountsService/users

[Install]
WantedBy=multi-user.target
